# ramdisk/rootfs/overlay/cv1812cp_licheerv_nano_sd/etc/rc.local
#!/bin/sh

echo Hello World
/opt/stress-test/stress-test-all.sh &
exit 0
