# buildroot/board/cvitek/SG200X/overlay/etc/profile.d/app_store.sh
export PATH=/maixapp/apps/app_store/app_store:$PATH
