# buildroot/board/cvitek/SG200X/overlay/etc/rc.local
#!/bin/sh

echo "HELLO WORLD"

if [ -e /boot/maixcam ]
then
	/maixapp/apps/launcher/launcher &
fi
